# backend/printdash-backend.service
[Unit]
Description=printdash FastAPI backend
After=network.target

[Service]
Type=simple
User=printdash
WorkingDirectory=/opt/printdash-backend/backend
EnvironmentFile=/etc/printdash/env
ExecStart=/opt/printdash-backend/backend/.venv/bin/uvicorn app.main:app --host 0.0.0.0 --port 8000
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
